# modules/agora-election/templates/supervisor_aelection.conf.erb
[program:aelection]
command=/home/aelection/start.sh
autostart=true
autorestart=true
user=aelection