# modules/agora-election/templates/supervisor_aelection.conf.erb
[program:aelection]
command=/home/aelection/agora-election/agora_election/app.py
directory=/home/aelection/agora-election/agora_election
autostart=true
autorestart=true
user=aelection
environment=PATH="/home/aelection/.virtualenvs/agora-election/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"